# debian/DEBIAN/postinst
#!/bin/sh
set -e

# Source debconf library.
. /usr/share/debconf/confmodule

db_get tknikabackups/backupsusername
username="$RET"
sed -i "s/www-data/$username/" /etc/cron.d/tknikabackups
if test "x$username" != xwww-data -a "x$username" != x
then
    if test x`grep $username /etc/passwd` = x
    then
        adduser --system --home /var/lib/tknikabackups --group $username
    fi
    # set rwx permissions for www-data and the backup user in the cache and log directories
    # as described in http://symfony.com/doc/current/book/installation.html#configuration-and-setup
    setfacl  -R -m u:www-data:rwx -m u:$username:rwx /var/cache/tknikabackups
    setfacl -dR -m u:www-data:rwx -m u:$username:rwx /var/cache/tknikabackups
    setfacl  -R -m u:www-data:rwx -m u:$username:rwx /var/log/tknikabackups
    setfacl -dR -m u:www-data:rwx -m u:$username:rwx /var/log/tknikabackups
    # the same in the spool directory
    setfacl  -R -m u:www-data:rwx -m u:$username:rwx /var/spool/tknikabackups
    setfacl -dR -m u:www-data:rwx -m u:$username:rwx /var/spool/tknikabackups
fi
chown $username.$username /var/cache/tknikabackups /var/log/tknikabackups

a2ensite tknikabackups
# reload apache
if [ -f /etc/init.d/apache2 ]
then
    invoke-rc.d apache2 reload
fi

db_stop
